# setup/splitflap.service
[Unit]
Description=Splitflap OS Web Server
After=network.target splitflap-network.service
Wants=splitflap-network.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/splitflap-os/server
ExecStart=/usr/bin/python3 -c "import app; app.app.run(host='0.0.0.0', port=80)"
Restart=on-failure
RestartSec=5
Environment=SPLITFLAP_CONFIG=/opt/splitflap-os/server/settings.json

[Install]
WantedBy=multi-user.target
